# debian/jenkins-slave.upstart.in
description "jenkins: Jenkins Slave Agent"
author "James Page <james.page@ubuntu.com>"

start on runlevel [2345]
stop on runlevel [!2345]

pre-start script
    [ -r /etc/default/jenkins-slave.default ] && . /etc/default/jenkins-slave
    [ -n "$JENKINS_URL" ] || { stop; exit 0 }
    $JENKINS_ROOT/bin/download-slave.sh
    mkdir $JENKINS_RUN  > /dev/null 2>&1  || true
    chown -R $JENKINS_USER $JENKINS_RUN || true
end script

script
    [ -r /etc/default/jenkins-slave.default ] && . /etc/default/jenkins-slave
    exec su -l $JENKINS_USER --shell=/bin/bash -c \
        "daemon --name=$NAME --inherit --output=$JENKINS_SLAVE_LOG \
         -- $JAVA $JAVA_ARGS -jar $JENKINS_RUN/slave.jar $JENKINS_ARGS" 
end script
